# files/etc/systemd/system/clash-verge-service.service
[Unit]
Description=Clash Verge Service
After=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/clash-verge-service
# 注：新机器上如需指定运行用户，取消注释并替换为实际用户名：
# User=YOUR_USERNAME
# Group=YOUR_USERNAME
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
