# files/etc/systemd/system/mihomo.service
[Unit]
Description=mihomo - Clash Meta Core (Transparent Proxy)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/mihomo -d /etc/mihomo
Restart=always
RestartSec=5
LimitNOFILE=1048576

[Install]
WantedBy=multi-user.target
