# files/etc/systemd/system/ppd-profile-monitor.service
[Unit]
Description=PPD Profile Power Limit Monitor

[Service]
Type=exec
ExecStart=/usr/local/bin/ppd-power-tune.sh daemon
Restart=on-failure
RestartSec=5
TimeoutStartSec=10

[Install]
WantedBy=multi-user.target
